# supervisor/tftp.ini
[program:tftp]
command=/usr/sbin/in.tftpd -4 -L -s -u nobody -P /tmp/tftpd.pid -v -p -B 1468 /tftpboot
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
